# scripts/com.ai-tips.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-tips.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPO_DIR/scripts/sync_ai_tips.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>HOME_DIR/Library/Logs/ai-tips-sync.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_DIR/Library/Logs/ai-tips-sync.log</string>
</dict>
</plist>
